# deploy/key49-fetch-api.service
[Unit]
Description=Key49-Fetch REST API Gateway
Documentation=https://github.com/grisbi-ia/key49fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/key49-fetch
EnvironmentFile=/opt/key49-fetch/.env
ExecStart=/opt/key49-fetch/.venv/bin/python -m src.api.app

Restart=on-failure
RestartSec=5

MemoryMax=256M
CPUQuota=100%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=key49-fetch-api

[Install]
WantedBy=multi-user.target
