# deploy/key49-fetch.service
[Unit]
Description=Key49-Fetch SRI Document Download Worker
Documentation=https://github.com/grisbi-ia/key49fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/key49-fetch
Environment="HOME=/root"
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/key49-fetch/.env
ExecStart=/opt/key49-fetch/.venv/bin/python -m src.orchestrator \
    --config /opt/key49-fetch/config/companies.json \
    --output /data/key49-fetch/xml_downloads

MemoryMax=1G
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=key49-fetch

[Install]
WantedBy=multi-user.target
